# printer_data/config/configFiles/rc.local
#!/bin/sh -e
# rc.local

echo 50 > /sys/class/backlight/10-0045/brightness
exit 0

